# etc/systemd/system/getty@tty1.service.d/skip-prompt.conf
[Service]
ExecStart=
ExecStart=-/sbin/agetty -o '-p -f -- \\u' --skip-login --nonewline --noissue --autologin dolphinos %I $TERM
Type=simple
Environment=XDG_SESSION_TYPE=x11
